# static/etc/services.d/jackett/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid bytesized mono /app/Jackett/JackettConsole.exe
